# etc/update-motd.d/REAME
#!/bin/sh
#
# This dynamic MOTD based on http://nickcharlton.net/posts/debian-ubuntu-dynamic-motd.html
